# override/etc/s6-overlay/s6-rc.d/quake2_data/up
#!/command/execlineb
foreground { /command/with-contenv /bin/bash /usr/local/bin/quake2-install }
