# override/etc/s6-overlay/s6-rc.d/quake2_data/up
#!/command/execlineb
# SPDX-FileCopyrightText: 2022 - 2024 Daniel Wolf <nephatrine@gmail.com>
#
# SPDX-License-Identifier: ISC

foreground { /command/with-contenv /bin/sh /usr/local/bin/quake2-install }
