# override/etc/s6-overlay/s6-rc.d/mount_setup/up
#!/command/execlineb

foreground { /bin/chown guardian:users /mnt/shared }
foreground { /bin/chown guardian:users /mnt/mirror }
